# build/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true

[program:spot_handler]
command=python manage.py spot_handler
directory=/app
autostart=true
autorestart=true
