# build/docker/supervisord.conf
[supervisord]
logfile=/dev/stdout
loglevel=info
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true

[program:spot_handler]
command=python manage.py spot_handler
directory=/app
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
